# ibkr-gateway/config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1024x768x24
autorestart=true
priority=1
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb-err.log

[program:ibgateway]
command=/scripts/start-ibgateway.sh
environment=DISPLAY=":1"
user=trader
autorestart=true
priority=10
startsecs=30
stdout_logfile=/var/log/ibgateway.log
stderr_logfile=/var/log/ibgateway-err.log

[program:fastapi]
command=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8080 --workers 1
directory=/app
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/fastapi.log
stderr_logfile=/var/log/fastapi-err.log
environment=PYTHONUNBUFFERED="1"
